# rootfs/etc/services.d/filebot-watcher/run
#!/usr/bin/with-contenv bash

# start wait loop
echo "**** waiting for instructions ****"

while true
do
  sleep 86400
done
